# deploy/systemd/nl2sql.service
[Unit]
Description=NL-to-SQL Backend
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/ubuntu/nl2sql
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down
TimeoutStartSec=300
User=ubuntu

[Install]
WantedBy=multi-user.target
